# home/private_dot_config/systemd/user/check-updates.timer
[Unit]
Description=Run checkup script periodically

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
